# prototype/main/im-config/fcitx5/fcitx5-chewing/asset/overlay/etc/profile.d/input-method.sh


##
## ## Arch Wiki
##
## * [Fcitx5](https://wiki.archlinux.org/title/Fcitx5)
## * [Fcitx5](https://wiki.archlinuxcn.org/zh-tw/Fcitx5)
##
##


GTK_IM_MODULE=fcitx
QT_IM_MODULE=fcitx
XMODIFIERS=@im=fcitx
SDL_IM_MODULE=fcitx
INPUT_METHOD=fcitx
GLFW_IM_MODULE=ibus
